# crontab
0 */2 * * * root python /src/app/run.py 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job

